# etc/supervisord.d/nginx-conf.ini
[program:nginx-conf]
command = /opt/qnib/nginx/bin/conf.sh
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
